# srcs/InfluxDB/supervisord.conf
[supervisord]
nodaemon=true


[program:influxdb]
command = ./usr/sbin/influxd
autorestart = true

[program:telegraf]
command = telegraf -config etc/telegraf.conf
autorestart = true