# srcs/InfluxDB/supervisord.conf
[supervisord]
nodaemon=true


[program:influxdb]
command = ./usr/sbin/influxd
autorestart = true
