# srcs/Mysql/supervisord.conf
[supervisord]
nodaemon=true
user=root


[program:mysql]
command:rc-service mariadb start
autorestart = true
