# srcs/Mysql/supervisord.conf
[supervisord]
nodaemon=true


[program:mysql]
command:rc-service mariadb start
