# srcs/nginx/supervisord.conf
[supervisord]
nodaemon=true


[program:nginx]
command = nginx -g "daemon off;"
autorestart = true

[program:telegraf]
command = telegraf -config etc/telegraf.conf
autorestart = true

[program:sshd]
command = /usr/sbin/sshd -D
autorestart = true
