# srcs/Ftps/supervisord.conf
[supervisord]
user=root
nodaemon=true
pidfile=/var/run/supervisord.pid
loglevel=debug

[program:vsftpd]
command= /usr/sbin/vsftpd /etc/vsftpd/vsftpd.conf
#autorestart = true

[program:telegraf]
command = telegraf -config etc/telegraf.conf
#autorestart = true

[eventlistener:processes]
command= /tmp/stop-supervisor.sh
events=PROCESS_STATE_STOPPED, PROCESS_STATE_EXITED, PROCESS_STATE_FATAL