# srcs/Ftps/supervisord.conf
[supervisord]
nodaemon=true
user=root


[program:vsftpd]
command= ./sbin/service vsftpd start
#command= rc-service vsftpd start
