# srcs/Ftps/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:vsftpd]
command= ./sbin/service vsftpd start
autorestart = true

[program:telegraf]
command = telegraf -config etc/telegraf.conf
autorestart = true

