# srcs/Ftps/supervisord.conf
[supervisord]
nodaemon=true


[program:vsftpd]
command:./sbin/service vsftpd start
